# ubuntu/maverick/applications/tde-guidance/debian/kde-guidance-trinity.postinst
#!/bin/sh

# Remove shutdown and reboot links; this init script does not need them.
rm -f /etc/rc0.d/K20kde-guidance /etc/rc6.d/K20kde-guidance

#remove script from older times
update-rc.d -f displayconfig-hwprobe.py remove

#DEBHELPER#
